# ops/launchd/chatbot.health-check.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>studio.norikane.chatbot.health-check</string>
  <key>WorkingDirectory</key>
  <string>/Users/norikene_satoshi/projects/norikane_satoshi_HP</string>
  <key>ProgramArguments</key>
  <array>
    <string>/opt/homebrew/bin/pnpm</string>
    <string>tsx</string>
    <string>scripts/chatbot/tier1-health-check.ts</string>
  </array>
  <key>StartCalendarInterval</key>
  <dict>
    <key>Hour</key>
    <integer>4</integer>
    <key>Minute</key>
    <integer>0</integer>
  </dict>
  <key>StandardOutPath</key>
  <string>/Users/norikene_satoshi/Library/Logs/norikane-chatbot-health-check.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/norikene_satoshi/Library/Logs/norikane-chatbot-health-check.err.log</string>
</dict>
</plist>
